# rpi-admin/copy-usb-service@.service
[Unit]
BindTo=%i.device
After=%i.device

[Service]
Type=oneshot
TimeoutStartSec=10
ExecStart=/root/copy_data.sh /%I
